# packaging/debian-control/prerm.in
#!/bin/sh

# errors shouldn't cause script to exit
set +e

initctl stop shiny-server 2>/dev/null

systemctl stop shiny-server.service 2>/dev/null
systemctl disable shiny-server.service 2>/dev/null
rm -rf /etc/systemd/system/shiny-server.service
systemctl daemon-reload

# clear error termination state
set -e
